# etc/systemd/swap.service
[Unit]
Description=Create Swapfile

[Service]
Type=oneshot
User=ubuntu
ExecStart=/home/ubuntu/scewpt/etc/bin/create_swapfile.sh

[Install]
WantedBy=multi-user.target